# assets/services/cue.service.in
[Unit]
Description=cue Reminder Daemon

[Service]
ExecStart={{BINARY_PATH}} serve
Restart=always

[Install]
WantedBy=default.target
